# supervisord/delayed_job.conf
[program:delayed_job]
user=openproject
directory=/var/www/openproject
command=/var/www/openproject/docker/scripts/env.sh bundle exec bin/delayed_job run --exit-on-complete
autorestart=true
